# etc/cvmfs-user-pub.cron
* * * * * root /usr/libexec/cvmfs-user-pub/snapshots
* * * * * root sleep 30; /usr/libexec/cvmfs-user-pub/snapshots
